# deploy/akces-ngrok.service
[Unit]
Description=Akces Hub - Ngrok Tunnels (Pro + Free)
After=akces-hub.service network-online.target
Wants=network-online.target
Requires=akces-hub.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/akces-hub
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/.local/bin/ngrok start --all --config /home/pi/.config/ngrok/ngrok.yml --config /home/pi/clean_v32/deploy/ngrok.yml --log=stdout
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
