# deploy/akces-ngrok.service
[Unit]
Description=Akces Hub - Ngrok Tunnel (watchdog)
After=akces-hub.service network-online.target
Wants=network-online.target
Requires=akces-hub.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/akces-hub
ExecStart=/bin/bash /home/pi/akces-hub/deploy/ngrok_watchdog.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
