# deploy/akces-backup.service
[Unit]
Description=Akces Hub - Daily Backup & Cloud Export
After=akces-hub.service

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/akces-hub
ExecStart=/usr/bin/python3 /home/pi/akces-hub/deploy/backup_cloud.py
TimeoutStartSec=600
